# glibc-2.3/debian/nscd/DEBIAN/postrm
#! /bin/sh
set -e

if [ "x$1" = xpurge ]; then
    update-rc.d nscd remove > /dev/null
fi

exit 0
