# glibc-2.3/debian/nscd/DEBIAN/postinst
#! /bin/sh
set -e

if [ "$1" = configure ]; then
    update-rc.d nscd defaults 19 > /dev/null
    /etc/init.d/nscd restart
fi
exit 0
