# freebsd-utils/debian/kfreebsd-common.postinst
#!/bin/sh
set -e
rm -f /etc/init.d/kfreebsd-common
update-rc.d kfreebsd-common remove >/dev/null || exit 0

#DEBHELPER#

exit 0
