# freebsd-utils/debian/kfreebsd-common.postinst
#!/bin/sh
set -e
update-rc.d kfreebsd-common remove >/dev/null || exit 0
rm -f /etc/init.d/kfreebsd-common
