# freebsd-utils/debian/kldutils.postinst
#!/bin/sh
set -e
rm -f /etc/init.d/kldutils
update-rc.d kldutils remove >/dev/null || exit 0
